# deployer.service.template
[Unit]
Description=Simple App Deployer
After=multi-user.target

[Service]
Type=idle
ExecStart=
Restart=always

[Install]
WantedBy=multi-user.target