# deployer.service.template
[Unit]
Description=Simple App Deployer
After=network.target

[Service]
Type=idle
ExecStart=/bin/bash /full/path/to/service.sh
Restart=always

[Install]
WantedBy=multi-user.target